# phpmyadmin/srcs/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=rc-service nginx start
stdout_logfile=/var/log/supervisord/nginx.log

[program:php-fpm7]
command=rc-service php-fpm7 start
stdout_logfile=/var/log/supervisord/php-fpm7.log
